# config/systemd/headscale.service.d/override.conf
[Service]
Environment=HEADSCALE_CONFIG=/etc/headscale/config.yaml
ReadWritePaths=/etc/headscale
