# scripts/systemd/dnsdist.service.d/10-no-port53.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dnsdist \
  --supervised \
  --disable-syslog \
  -C /etc/dnsdist/dnsdist.conf \
  --local=127.0.0.1:5300 \
  --local=[::1]:5300
