# archive/mk/attic/attic.service
[Unit]
Description=Attic CAS (LAN binary cache)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Ensure log directory exists before attempting file-based logging
ExecStartPre=/bin/mkdir -p /volume1/homelab/attic/logs

# Fail loudly if the expected binary is missing or not executable
ExecStartPre=/bin/test -x /usr/local/bin/atticd

# Start Attic with an explicit, absolute config path
ExecStart=/usr/local/bin/atticd --config /volume1/homelab/attic/config.toml

# Set a stable working directory for predictable behavior
WorkingDirectory=/volume1/homelab/attic

# Runs as root intentionally to access NAS storage paths
User=root
Group=root

# Ensure predictable permissions
UMask=0022

TimeoutStartSec=15
TimeoutStopSec=30

ExecReload=/bin/false
LimitNOFILE=1048576

# File logging (no journald)
StandardOutput=append:/volume1/homelab/attic/logs/attic.log
StandardError=append:/volume1/homelab/attic/logs/attic.log

Restart=on-failure
RestartSec=3

# Hardening
ProtectSystem=strict
ReadWritePaths=/volume1/homelab/attic

NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes

AmbientCapabilities=
CapabilityBoundingSet=

KillMode=process

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
