# config/systemd/homelab-nft-rollback.timer
[Unit]
Description=Rollback homelab nftables if not confirmed

[Timer]
OnActiveSec=120

[Install]
WantedBy=timers.target
