# config/systemd/tailscaled-guest.service
# tailscaled-guest.service
[Unit]
Description=Tailscale guest role (exit node only)
After=network.target tailscaled.service
Requires=tailscaled.service

[Service]
Type=oneshot
ExecStart=/usr/bin/tailscale up --reset \
  --advertise-exit-node \
  --accept-dns=true \
  --accept-routes=false
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
