# config/systemd/unbound-ctl-fix.path
[Unit]
Description=Watch for Unbound control socket and trigger fix

[Path]
PathExists=/run/unbound.ctl

[Install]
WantedBy=multi-user.target
