# config/systemd/unbound.service
# /home/julie/src/homelab/config/systemd/unbound.service
[Unit]
Description=Unbound DNS resolver
Documentation=man:unbound(8)
After=network.target

[Service]
User=root
Group=root
Type=simple
# Unbound runs in foreground; it will drop to "unbound" per unbound.conf
ExecStart=/usr/sbin/unbound -d -c /etc/unbound/unbound.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
LimitNOFILE=400000

# Capabilities: allow binding to privileged ports; do not block setuid/setgid
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_CHOWN
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_CHOWN

[Install]
WantedBy=multi-user.target
