# 10.89.12.4/etc/systemd/system/coredns.service
# /etc/systemd/system/coredns.service
[Unit]
Description=CoreDNS
After=network.target

[Service]
User=coredns
Group=coredns
ExecStart=/usr/local/bin/coredns -conf /etc/coredns/Corefile
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
# After modification
# sudo systemctl daemon-reload
# sudo systemctl reset-failed
# sudo systemctl enable --now coredns
# sudo journalctl -u coredns -f
